# web/files/etc/cron.d/magento2
* * * * * www-data /usr/local/bin/php /var/www/magento-formation/src/bin/magento cron:run | grep -v "Ran jobs by schedule" >> /var/www/magento-formation/src/var/log/magento.cron.log
